# deploy/xenocosm.service
[Unit]
Description=xenocosm

[Service]
ExecStart=/usr/local/bin/xenocosm

[Install]
WantedBy=multi-user.target
